# systemd/shadow-reducer.service
[Unit]
Description=Shadow Creatures Reducer (Kismet -> ghost_state.json)
After=network.target

[Service]
Type=simple
User=jermarti
WorkingDirectory=/home/jermarti/maind-deploy
EnvironmentFile=/home/jermarti/maind-deploy/config/shadow.env
ExecStart=/usr/bin/python3 /home/jermarti/maind-deploy/src/reducer/reducer.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
