# systemd/shadow-backend.service
[Unit]
Description=Shadow Creatures Backend (Exhibition: WebSocket + static files)
After=network.target

[Service]
Type=simple
User=jermarti
WorkingDirectory=__DEPLOY_DIR__
EnvironmentFile=__DEPLOY_DIR__/config/shadow.env
ExecStart=/usr/bin/python3 __DEPLOY_DIR__/src/backend/server.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
