# systemd/shadow-kiosk.service
[Unit]
Description=Shadow Creatures Kiosk (Chromium fullscreen)
After=shadow-backend.service
Wants=shadow-backend.service

[Service]
Type=simple
# The kiosk script waits for the X11/Wayland display to be ready
ExecStart=/home/jermarti/maind-deploy/scripts/kiosk.sh http://localhost:8080
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
