# systemd/shadow-power.service
[Unit]
Description=Shadow Creatures Power Monitor (X1201 UPS)
# EXHIBITION MODE: This service is NOT needed — no battery/UPS hardware.
# Do not enable on exhibition deployments.
After=multi-user.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 __DEPLOY_DIR__/scripts/power_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
